# arch-edition/porco.service
[Unit]
Description=Porco Music Bot (Arch edition — motor mpv)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/porco-music-bot/arch-edition/engine.py
Restart=always
RestartSec=5
StandardOutput=append:%h/porco-music-bot/bot.log
StandardError=inherit
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
